# roles/minio/templates/minio.service.j2
[Unit]
Description=MinIO Object Storage
Documentation=https://min.io/docs/minio/linux/index.html
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User={{ minio_user }}
Group={{ minio_group }}
EnvironmentFile=/etc/default/minio
ExecStart={{ minio_install_dir }}/minio server $MINIO_VOLUMES $MINIO_OPTS
Restart=on-failure
RestartSec=10
LimitNOFILE=65536
TasksMax=infinity
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
